# templates/ckb-bench-miner.service.j2
[Unit]
Description=CKB Bench Miner
After=network.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=no
RestartSec=2
User=ckb
Group=ckb
Environment="RUST_LOG=info,ckb-bench=debug"
ExecStart={{ ckb_benchmark_workspace }}/ckb-bench miner --rpc-urls {{ ckb_benchmark_rpc_urls[0] }} --mining-interval-ms 1000

[Install]
WantedBy=multi-user.target
